# projectpi_releasecandidate/SOFTWARE/rc.local
#!/bin/sh -e

ifup wlan0

projectpi &

exit 0